# rootfs/usr/lib/systemd/system/pam_login_ng.service
[Unit]
Description=pam service for login_ng
Before=greetd.service

[Service]
Type=dbus
BusName=org.zbus.login_ng
ExecStart=/usr/bin/pam_login_ng-service
Restart=always
IgnoreSIGPIPE=no

[Install]
WantedBy=multi-user.target