# rootfs/usr/lib/systemd/system/greetd.service.d/override.conf
[Unit]
After=pam_login_ng.service
Wants=pam_login_ng.service
