# rootfs/usr/lib/systemd/system/greetd.service.d/override.conf
[Unit]
After=pam_polyauth.service
Wants=pam_polyauth.service
